# scripts/minecraft-god.service
[Unit]
Description=Minecraft God Server (BDS + LLM Backend)
After=network.target

[Service]
Type=forking
User=drew
WorkingDirectory=/home/drew/projects/minecraft-god
ExecStart=/home/drew/projects/minecraft-god/scripts/start.sh
ExecStop=/home/drew/projects/minecraft-god/scripts/stop.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
