# scripts/minecraft-god.service
[Unit]
Description=Minecraft God Server (Paper MC + LLM Backend)
After=network.target

[Service]
Type=forking
# Update these paths for your installation
User=%i
WorkingDirectory=/path/to/minecraft-god
ExecStart=/path/to/minecraft-god/scripts/start.sh
ExecStop=/path/to/minecraft-god/scripts/stop.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
